# deployment/crontab
# Fuel Monitor Cron Schedule
# Runs every hour at the top of the hour
# Logs to /var/log/fuel-monitor/cron.log

# Load environment variables from .env and run monitor-fuel
0 * * * * cd /app && export $(grep -v '^#' .env | grep -v '^$' | xargs) && npx tsx monitor-fuel.ts >> /var/log/fuel-monitor/cron.log 2>&1

# Keep a newline at the end for cron to work properly
